# .config/systemd/user/letsencrypt.timer
[Unit]
Description=Update LetsEncrypt certificates once a month

[Timer]
OnCalendar=monthly

[Install]
WantedBy=timers.target

